# .config/autostart/polybar.desktop
[Desktop Entry]
Exec=polybar
Icon=
Name=polybar
Path=
Terminal=False
Type=Application
