# .config/autostart/syncthing.desktop
[Desktop Entry]
Exec=syncthing
Icon=
Name=syncthing
Path=
Terminal=False
Type=Application
